# .config/autostart/mate-panel.desktop
[Desktop Entry]
Exec=/usr/bin/mate-panel
Name=MATE Panel
Comment=The Panel for the MATE Desktop Environment
Type=Application
